# panoptic.service
[Unit]
Description=panoptic session monitor
After=network.target

[Service]
Type=simple
WorkingDirectory=__INSTALL_DIR__
ExecStart=__PYTHON3__ panoptic_cli.py serve
Restart=on-failure
RestartSec=5

# Logging goes to journald by default; use journalctl -u panoptic to view.
# Optionally redirect to files like the launchd configuration:
# StandardOutput=append:__INSTALL_DIR__/logs/stdout.log
# StandardError=append:__INSTALL_DIR__/logs/stderr.log

[Install]
WantedBy=default.target
